# config/etc/systemd/system/alsa_vol_from_usb_hid.service
[Unit]
Description=ALSA Volume Control from USB HID Consumer Control

[Service]
Type=simple
User=pi
ExecStart=%TARGET%
Restart=on-failure
